# legacy/deploy/systemd/homehub-certbot-renew.timer
[Unit]
Description=Check HomeHub public IP certificate twice daily

[Timer]
OnCalendar=*-*-* 00,12:20:00
RandomizedDelaySec=30m
Persistent=true
Unit=homehub-certbot-renew.service

[Install]
WantedBy=timers.target
